# cluster/roles/vip/templates/check-keepalived-status.timer
[Unit]
Description=Check keepalived status every 15 seconds

[Timer]
OnBootSec=1min
OnUnitActiveSec=15s

[Install]
WantedBy=timers.target
